# tools/vm-service/vm-mcp.service
[Unit]
Description=VM MCP Server
After=vm.service

[Service]
Type=simple
Environment=VM_SSH_PORT=2222
Environment=VM_SSH_USER=joao
Environment=MCP_PORT=8765
ExecStart=node %h/proj/angst/tools/vm-mcp/dist/index.js
WorkingDirectory=%h/proj/angst/tools/vm-mcp
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
