# tools/vm-service/vm.service
[Unit]
Description=NixOS VM (angst)
After=network.target

[Service]
Type=simple
Environment=NIX_DISK_IMAGE=%h/proj/angst/personal.qcow2
Environment=ANGST_REPO=%h/proj/angst
Environment=QEMU_NET_OPTS=hostfwd=tcp::2222-:22
ExecStart=%h/proj/angst/result/bin/run-personal-vm
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
